# hardware_specific/odroid-xu4/update-fans.timer
[Unit]
Description=Update the system fan and temperature setting at startup

[Timer]
OnStartupSec=5

[Install]
WantedBy=multi-user.target
