# hardware_specific/odroid-xu4/update-fans.service
[Unit]
Description=update-fans

[Service]
Type=simple
ExecStart=/usr/local/bin/update-fans.sh
